# loopback@.service
[Unit]
Description=LoopbackStart
[Service]
ExecStart=/bin/bash -c "losetup -f /home/jarda/%I"
[Install]
WantedBy=multi-user.target
